# profile.d/borealis-pg-export-autossh-path.sh
#!/usr/bin/env sh

export PATH="$HOME/.borealis-pg/autossh:$PATH"
